# deploy/systemd/genie-governor.service
[Unit]
Description=GeniePod Memory Governor
Documentation=https://github.com/GeniePod/genie-claw
After=network.target
Wants=genie-core.service

[Service]
Type=notify
ExecStart=/opt/geniepod/bin/genie-governor
Environment=GENIEPOD_CONFIG=/etc/geniepod/geniepod.toml
Environment=RUST_LOG=info
Restart=always
RestartSec=3
WatchdogSec=30
MemoryMax=32M
CPUWeight=50

# Pin to the system bucket (issue #25) — governor only polls memory and writes
# drop-in configs, keep it off the inference cores (2-5).
CPUAffinity=0 1

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/geniepod/data /run/geniepod /etc/systemd/system/genie-llm.service.d /etc/systemd/system/genie-ai-runtime.service.d
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
