# deploy/systemd/genie-mqtt.service
[Unit]
Description=GeniePod MQTT Broker (Mosquitto)
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/mosquitto -c /etc/geniepod/mosquitto.conf
Restart=always
RestartSec=3
MemoryMax=16M

# Pin to the system bucket (issue #25) — broker traffic is light and event-
# driven, keep it off the inference cores (2-5).
CPUAffinity=0 1

ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/geniepod/data/mqtt
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=geniepod.target
