# deploy/systemd/genie-api.service
[Unit]
Description=GeniePod System Dashboard
Documentation=https://github.com/GeniePod/genie-claw
After=network.target genie-governor.service

[Service]
Type=simple
ExecStart=/opt/geniepod/bin/genie-api
Environment=GENIEPOD_CONFIG=/etc/geniepod/geniepod.toml
Environment=RUST_LOG=info
Restart=always
RestartSec=5
MemoryMax=32M

# Pin to the system bucket (issue #25) so the dashboard's HTTP loop doesn't
# steal cycles from whisper (2-3), llama-server (4), or genie-core (5).
CPUAffinity=0 1

ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/geniepod/data
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=geniepod.target
