# deploy/systemd/genie-health.service
[Unit]
Description=GeniePod Health Monitor
Documentation=https://github.com/GeniePod/genie-claw
After=network.target genie-core.service

[Service]
Type=notify
ExecStart=/opt/geniepod/bin/genie-health
Environment=GENIEPOD_CONFIG=/etc/geniepod/geniepod.toml
Environment=RUST_LOG=info
Restart=always
RestartSec=5
WatchdogSec=60
MemoryMax=32M
CPUWeight=20

# Pin to the system bucket (issue #25) — health monitor is a periodic poller,
# keep it off the inference cores (2-5).
CPUAffinity=0 1

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/geniepod/data
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
